# ansible/roles/observability/files/docker-state-exporter.timer
[Unit]
Description=Docker state exporter timer (managed by Ansible)

[Timer]
# Prometheus scrapes every 30s; emit fresh textfile every 15s so a scrape
# never sees a sample older than one cycle.
OnBootSec=10s
OnUnitActiveSec=15s
AccuracySec=1s
Unit=docker-state-exporter.service

[Install]
WantedBy=timers.target
