# ansible/roles/wireguard/templates/wg-watchdog.service.j2
[Unit]
Description=WireGuard handshake watchdog (managed by Ansible)
After=wg-quick@wg0.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/wg-watchdog.sh
Nice=19
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
